# systemd/moonbit-clean.service
[Unit]
Description=MoonBit Automated Cleaning
Documentation=https://github.com/Nomadcxx/moonbit
After=network.target
# Ensure scan has run before cleaning
After=moonbit-scan.service

[Service]
Type=oneshot
# Run with --force to actually clean (not dry-run)
ExecStart=/usr/local/bin/moonbit clean --force
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moonbit-clean

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/cache /var/tmp /var/log /root/.cache /root/.local/share

[Install]
WantedBy=multi-user.target
