# systemd/moonbit-scan.service
[Unit]
Description=MoonBit System Scan
Documentation=https://github.com/Nomadcxx/moonbit
After=network.target

[Service]
Type=oneshot
Environment=HOME=/root
Environment=XDG_CONFIG_HOME=/var/lib/moonbit/config
Environment=XDG_CACHE_HOME=/var/cache
Environment=XDG_DATA_HOME=/var/lib
StateDirectory=moonbit
CacheDirectory=moonbit
ExecStart=/usr/local/bin/moonbit scan --mode quick --no-prompt
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moonbit-scan

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/cache/moonbit /var/lib/moonbit

[Install]
WantedBy=multi-user.target
