# systemd/moonbit-scan.timer
[Unit]
Description=MoonBit Daily Scan Timer
Documentation=https://github.com/Nomadcxx/moonbit

[Timer]
# Run daily at 2 AM
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
